# multisig.notional.ventures.service
[Unit]
Description=Multisig Deploy Service

[Service]
Restart=always
RestartSec=3
ExecStart=cd /root/multisig.notional.ventures && $(which npm) start

[Install]
WantedBy=multi-user.target




