# multisig.notional.ventures.service
[Unit]
Description=Multisig Deploy Service

[Service]
Restart=always
RestartSec=3
WorkingDirectory=/root/multisig.notional.ventures
ExecStart=$(which npm) install && $(which npm) run build && $(which npm) start

[Install]
WantedBy=multi-user.target




